# deploy/drummajor-portal.service
# systemd unit for the Drum Major Portal (native, non-dockerized).
# Install to: /etc/systemd/system/drummajor-portal.service
#
#   sudo cp deploy/drummajor-portal.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now drummajor-portal
#
# On every (re)start this runs deploy/update-and-run.sh, which pulls new
# commits, rebuilds when needed, then execs `next start`.
# App lives at /home/musicserver/drummajor-portal, owned by `musicserver`.

[Unit]
Description=Drum Major Portal (Next.js)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=musicserver
Group=musicserver
WorkingDirectory=/home/musicserver/drummajor-portal
# Secrets + DB URL live here (chmod 640, owner root:musicserver).
EnvironmentFile=/etc/drummajor-portal/portal.env
# systemd ships a minimal PATH; add where node/npm/git/npx live so the
# update script can find them. Adjust if you use nvm or a custom prefix.
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/home/musicserver/drummajor-portal/deploy/update-and-run.sh
# Build can take a while on first run / big updates — don't time out the start.
TimeoutStartSec=0
Restart=on-failure
RestartSec=5
# Hardening (ProtectHome must stay false: the app lives under /home).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
