# deploy/drummajor-portal.service
# systemd unit for the Drum Major Portal (native, non-dockerized).
# Install to: /etc/systemd/system/drummajor-portal.service
#
#   sudo cp deploy/drummajor-portal.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now drummajor-portal
#
# Assumes the app is checked out at /opt/drummajor-portal, built with
# `npm run build`, and that Node 22 is available (see deploy/README.md).

[Unit]
Description=Drum Major Portal (Next.js)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=portal
Group=portal
WorkingDirectory=/opt/drummajor-portal
# Secrets + DB URL live here, root-readable only (chmod 600, owner root:portal).
EnvironmentFile=/etc/drummajor-portal/portal.env
# Adjust the node path to wherever Node 22 is installed on the server.
ExecStart=/usr/bin/node /opt/drummajor-portal/node_modules/.bin/next start -p 3000
Restart=on-failure
RestartSec=5
# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
